# templates/node_basic_server.service.erb
[Unit]
Description=no description given

[Service]
Type=simple
User=root
Group=root
# Load env vars from /etc/default/ and /etc/sysconfig/ if they exist.
# Prefixing the path with '-' makes it try to load, but if the file doesn't
# exist, it continues onward.
EnvironmentFile=-/etc/default/nodejs
EnvironmentFile=-/etc/sysconfig/nodejs
ExecStart=/opt/node-v6.9.1-linux-x64/bin/node /srv/node/server_<%= @name %>.js
Restart=always
WorkingDirectory=/

[Install]
WantedBy=multi-user.target
